# chef/cookbooks/crowbar/files/default/crowbar-sessions-sweep.cron
#!/bin/bash
su -s /bin/sh - crowbar sh -c "cd /opt/dell/crowbar_framework && RAILS_ENV=production ./bin/rake sessions:sweep" || :
