# chef/cookbooks/provisioner/templates/default/tftp.service.erb
[Unit]
Description=Tftp Server

[Service]
Type=simple
ExecStart=/usr/sbin/in.tftpd -u tftp -s <%= @tftproot %> -m /etc/tftpd.conf -L -a <%= @admin_ip %> -B 1024 -v
